# install/duktig.service
[Unit]
Description=Autostart service for Duktig
After=duktig-updater.service

[Service]
WorkingDirectory=/home/seard/Desktop/duktig-proj/duktig/
ExecStart=sudo node server.js
TimeoutStartSec=5
#Restart=always
StandardOutput=journal
StandardError=journal
SyslogIdentifier=duktig
User=root
Group=root

[Install]
WantedBy=multi-user.target
